# scripts/templates/devassist-orchestrator.service.j2
[Unit]
Description=developer-assistant orchestrator Hermes runtime (Telegram gateway)
PartOf=devassist.target
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=devassist
Group=devassist
WorkingDirectory=/srv/devassist/runtimes/orchestrator
Environment=HERMES_HOME=/srv/devassist/runtimes/orchestrator/.hermes
Environment=HOME=/srv/devassist/runtimes/orchestrator
Environment=HERMES_DEVASSIST_ROLE=orchestrator
Environment=DEVASSIST_HEALTH_PORT=8181
Environment=XDG_STATE_HOME=/srv/devassist/runtimes/orchestrator/.local/state
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
ExecStart=/usr/local/bin/devassist-orchestrator-runner
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/srv/devassist
ReadWritePaths=/srv/devassist/runtimes/orchestrator /srv/devassist/state /srv/devassist/logs
BindReadOnlyPaths=/srv/devassist/repo /srv/devassist/shared-skills /srv/devassist/shared-plugins /srv/devassist/releases

[Install]
WantedBy=devassist.target
