# scripts/templates/devassist.target.j2
[Unit]
Description=developer-assistant umbrella target
Wants=devassist-orchestrator.service devassist-planner.service devassist-architect.service devassist-executor.service devassist-reviewer.service
After=network-online.target

[Install]
WantedBy=multi-user.target
