# scripts/templates/devassist-reviewer.service.j2
[Unit]
Description=developer-assistant reviewer Hermes runtime
PartOf=devassist.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devassist
Group=devassist
SupplementaryGroups=docker
WorkingDirectory=/srv/devassist/runtimes/reviewer
Environment=HERMES_HOME=/srv/devassist/runtimes/reviewer/.hermes
Environment=HERMES_DEVASSIST_ROLE=reviewer
Environment=DEVASSIST_HEALTH_PORT=8185
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
ExecStart=/usr/local/bin/devassist-worker-runner
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=5

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/srv/devassist
ReadWritePaths=/srv/devassist/runtimes/reviewer /srv/devassist/state /srv/devassist/logs
BindReadOnlyPaths=/srv/devassist/repo /srv/devassist/shared-skills /srv/devassist/shared-plugins /srv/devassist/releases

[Install]
WantedBy=devassist.target
