# scripts/templates/devassist.target.j2
[Unit]
Description=developer-assistant umbrella target
Wants=omniroute.service devassist-orchestrator.service devassist-planner.service devassist-architect.service devassist-executor.service devassist-reviewer.service devassist-web.service
After=network-online.target

[Install]
WantedBy=multi-user.target
