# scripts/templates/devassist-architect.service.j2
[Unit]
Description=developer-assistant architect Hermes runtime
PartOf=devassist.target
After=network-online.target omniroute.service
Wants=network-online.target

[Service]
Type=simple
User=devassist
Group=devassist
WorkingDirectory=/srv/devassist/runtimes/architect
Environment=HERMES_HOME=/srv/devassist/runtimes/architect/.hermes
Environment=HERMES_DEVASSIST_ROLE=architect
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
ExecStart=/usr/local/bin/hermes run
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=5

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/srv/devassist
ReadWritePaths=/srv/devassist/runtimes/architect /srv/devassist/state /srv/devassist/logs
BindReadOnlyPaths=/srv/devassist/repo /srv/devassist/shared-skills /srv/devassist/shared-plugins /srv/devassist/releases

[Install]
WantedBy=devassist.target
