# scripts/templates/devassist-planner.service.j2
[Unit]
Description=developer-assistant planner Hermes runtime
PartOf=devassist.target
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=devassist
Group=devassist
WorkingDirectory=/srv/devassist/runtimes/planner
Environment=HERMES_HOME=/srv/devassist/runtimes/planner/.hermes
Environment=HOME=/srv/devassist/runtimes/planner
Environment=HERMES_DEVASSIST_ROLE=planner
Environment=DEVASSIST_HEALTH_PORT=8182
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
ExecStart=/usr/local/bin/devassist-worker-runner
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/srv/devassist
ReadWritePaths=/srv/devassist/runtimes/planner /srv/devassist/state /srv/devassist/logs
BindReadOnlyPaths=/srv/devassist/repo /srv/devassist/shared-skills /srv/devassist/shared-plugins /srv/devassist/releases

[Install]
WantedBy=devassist.target
