# scripts/templates/devassist-executor.service.j2
[Unit]
Description=developer-assistant executor Hermes runtime
PartOf=devassist.target
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=devassist
Group=devassist
SupplementaryGroups=docker
WorkingDirectory=/srv/devassist/runtimes/executor
Environment=HERMES_HOME=/srv/devassist/runtimes/executor/.hermes
Environment=HOME=/srv/devassist/runtimes/executor
Environment=HERMES_DEVASSIST_ROLE=executor
Environment=DEVASSIST_HEALTH_PORT=8184
Environment=PYTHONPATH=/srv/devassist/repo/src
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
# AUDIT-001 / TKT-033: invoke runtime_check.check_runtime() blockingly before
# the unit's ExecStart so a failed invariant aborts boot with EX_CONFIG=78
# (developer_assistant.runtime_check.RUNTIME_CHECK_ABORT_EXIT_CODE) and emits
# the structured RUNTIME_CHECK_FAILED:<role>:<invariant_name> marker on stderr
# (captured by journald). RestartPreventExitStatus=78 below ensures invariant
# aborts surface as a `failed` unit and never auto-restart.
ExecStartPre=/usr/bin/python3 -m developer_assistant.runtime_check
ExecStart=/usr/local/bin/devassist-worker-runner
Restart=always
RestartSec=10s
RestartPreventExitStatus=78

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/srv/devassist
ReadWritePaths=/srv/devassist/runtimes/executor /srv/devassist/state /srv/devassist/logs
BindReadOnlyPaths=/srv/devassist/repo /srv/devassist/shared-skills /srv/devassist/shared-plugins /srv/devassist/releases

[Install]
WantedBy=devassist.target
