# scripts/templates/omniroute.service.j2
[Unit]
Description=OmniRoute v3.7.x routing layer for developer-assistant
PartOf=devassist.target
After=network-online.target
Wants=network-online.target
Before=devassist-orchestrator.service devassist-planner.service devassist-architect.service devassist-executor.service devassist-reviewer.service

[Service]
Type=simple
User=omniroute
Group=omniroute
WorkingDirectory=/srv/devassist/omniroute
EnvironmentFile=/srv/devassist/secrets/SELF-DEPLOY.env
ExecStart=/opt/omniroute/bin/omniroute serve --port 20128 --bind 127.0.0.1 --state-db /srv/devassist/omniroute/omniroute.db
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=5

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/opt/omniroute
ReadWritePaths=/srv/devassist/omniroute

[Install]
WantedBy=devassist.target
